# packages/folder-sync/photoframe-sync.timer
[Unit]
Description=Run Photoframe Sync periodically

[Timer]
OnUnitActiveSec=30
AccuracySec=1s
Persistent=true

[Install]
WantedBy=timers.target
