# packages/folder-sync/photoframe-sync@.service
[Unit]
Description=Photoframe Sync Service
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=user
Group=user
# The "%i" will be filled with instance name (from .path)
ExecStart=/usr/bin/env python3 /path/to/folder-sync.py -c /path/to/sync-config.yml -s %i
StandardOutput=journal
StandardError=journal
