# packages/folder-sync/photoframe-sync.timer
[Unit]
Description=Run Photoframe Sync periodically

[Timer]
OnUnitActiveSec=30s
AccuracySec=1s
Persistent=true
Unit=photoframe-sync@timer.service

[Install]
WantedBy=timers.target
